# hobday/crontab
0 12 * * * /usr/local/bin/python /root/main.py 
